# nzedb/s6/tmux/run
#!/bin/sh
cd /var/www/nZEDb/misc/update/nix/tmux
exec php start.php
