# nzedb/s6/tmux/run
#!/bin/sh
cd /var/www/nZEDb/misc/update/nix/tmux
php start.php
# Block indefinitely
exec tail -f /proc/`pidof tmux`/fd/1
