# nzedb/s6/sphinx/run
#!/bin/sh
exec searchd --nodetach -c /var/www/nZEDb/misc/sphinxsearch/sphinx.conf
